# kawagashira_nobuyuki/tyama_codeiq294.cron
30 23 * * * du -d 1 /home > /var/lib/bak/du-`date +\%Y\%m\%d`.log